# deploy/all-in-one/supervisord.conf
# Loaded via /etc/supervisor/supervisord.conf's include directive.
# No [supervisord] block here — main conf owns daemon settings.

[program:postgres]
command=/usr/lib/postgresql/16/bin/postgres -D /data/pgsql -c config_file=/data/pgsql/postgresql.conf
user=postgres
priority=10
autostart=true
autorestart=true
stdout_logfile=/var/log/akb/postgres.log
stderr_logfile=/var/log/akb/postgres.err
stopsignal=INT

[program:redis]
command=/usr/bin/redis-server --bind 127.0.0.1 --port 6379 --dir /var/lib/redis --save "" --appendonly no
user=redis
priority=15
autostart=true
autorestart=true
stdout_logfile=/var/log/akb/redis.log
stderr_logfile=/var/log/akb/redis.err

[program:minio]
command=/usr/local/bin/minio server /data/minio --address :9000 --console-address :9001
environment=MINIO_ROOT_USER="%(ENV_S3_ACCESS_KEY)s",MINIO_ROOT_PASSWORD="%(ENV_S3_SECRET_KEY)s"
user=root
priority=15
autostart=true
autorestart=true
stdout_logfile=/var/log/akb/minio.log
stderr_logfile=/var/log/akb/minio.err

[program:backend]
command=/opt/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
directory=/app/backend
user=root
priority=30
autostart=true
autorestart=true
stdout_logfile=/var/log/akb/backend.log
stderr_logfile=/var/log/akb/backend.err
stopwaitsecs=15

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
user=root
priority=40
autostart=true
autorestart=true
stdout_logfile=/var/log/akb/nginx.log
stderr_logfile=/var/log/akb/nginx.err

# One-shot seeder — waits for backend /readyz, then upserts demo user +
# vault + fixed PAT. Idempotent across restarts.
[program:seed]
command=/opt/venv/bin/python /usr/local/bin/akb-seed.py
user=root
priority=50
autostart=true
autorestart=false
startsecs=0
startretries=0
exitcodes=0
stdout_logfile=/var/log/akb/seed.log
stderr_logfile=/var/log/akb/seed.err
